# dnf-nag.service
[Unit]
Description=Nag the user if there are updates

[Service]
ExecStart=/usr/bin/dnf-nag.sh
